# .config/systemd/user/kakoune.service
[Unit]
Description=Kakoune

[Service]
ExecStart=/bin/zsh /home/user/.config/kak/script.sh

[Install]
WantedBy=default.target
